# deploy/vps/orchestrator-bot.service
[Unit]
Description=RivetCEO Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/Agent-Factory
Environment=PATH=/root/.local/bin:/usr/local/bin:/usr/bin
EnvironmentFile=/opt/Agent-Factory/.env
ExecStart=/root/.local/bin/poetry run python -m agent_factory.integrations.telegram.orchestrator_bot
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
