# deploy/vps/rivet-worker.service
[Unit]
Description=RIVET KB Ingestion Worker (24/7 Background Daemon)
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/Agent-Factory
Environment=PATH=/root/.local/bin:/usr/local/bin:/usr/bin
EnvironmentFile=/root/Agent-Factory/.env
ExecStart=/root/.local/bin/poetry run python scripts/rivet_worker.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Graceful shutdown
TimeoutStopSec=30
KillMode=process
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
