# deploy/vps/rivet-scheduler.timer
[Unit]
Description=Daily KB Building Scheduler Timer
Requires=rivet-scheduler.service

[Timer]
# Run every 4 hours (6 times daily)
OnCalendar=0/4:00:00
# Persistent: Run missed timers on boot
Persistent=true

[Install]
WantedBy=timers.target
